# AutomatedVehicle/Debug/AutomatedVehicle.service
[Unit]
Description=Automated Vehicle Service


[Service]
Type=oneshot

ExecStart=/home/pi/AutomatedVehicle

[Install]
WantedBy=multi-user.target
